# .claude-events/systemd/inotify-watch.service
[Unit]
Description=Inotify watch on .claude/ directory — cross-CLI event bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/srv/monorepo/.claude-events/inotify-watch.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

Environment="HOME=/home/will"
Environment="EVENT_DIR=/srv/monorepo/.claude-events"
Environment="CLAUDE_DIR=/srv/monorepo"

StandardOutput=journal
StandardError=journal
SyslogIdentifier=inotify-watch

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/
ReadWritePaths=/srv/monorepo/.claude-events/

MemoryMax=100M
TasksMax=4

[Install]
WantedBy=default.target
