# .claude-events/systemd/trigger-bridge.service
[Unit]
Description=Trigger bridge — poll state.json and dispatch to vibe-kit
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/srv/monorepo/.claude-events/trigger-bridge.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

Environment="HOME=/home/will"
Environment="EVENT_DIR=/srv/monorepo/.claude-events"
Environment="CLAUDE_DIR=/srv/monorepo"

StandardOutput=journal
StandardError=journal
SyslogIdentifier=trigger-bridge

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/
ReadWritePaths=/srv/monorepo/.claude-events/

MemoryMax=100M
TasksMax=4

[Install]
WantedBy=default.target
