# systemd/trieve.service
[Unit]
Description=Trieve Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose -f /srv/monorepo/docker-compose.trieve.yml up
ExecStop=/usr/bin/docker compose -f /srv/monorepo/docker-compose.trieve.yml down
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target